# examples/org-taube.service
[Unit]
Description=Run org-taube

[Service]
Type=oneshot
ExecStart=org-taube --auto
