# examples/org-taube.cron
# Sync mail and run org-taube every 5 minutes.
# Install with: crontab -e
*/5 * * * * mbsync -a && org-taube --auto
