# examples/org-taube.cron
# Run org-taube every 5 minutes.
# Install with: crontab -e
*/5 * * * * org-taube --auto
